# bin/install_archlinux/etc/crontab
*/30 * * * * /usr/bin/sh ${HOME}/bin/email/getmail.sh
